# s6/services.d/jikan-queue-2/finish
#!/usr/bin/execlineb -S0

echo Jikan Job Queue #2 terminated

# Terminate container
s6-svscanctl -t /var/run/s6/services
